# services/update_clock_gps.service
[Unit]
Description=Update clock with chronyc

[Service]
ExecStart=chronyc -a makestep
User=root
SyslogIdentifier=update_clock_gps
